# data/templates/gitea_openrc.ept
#!/sbin/openrc-run

name=gitea
command="/usr/bin/gitea"
command_args="web -c /etc/yunorock/<%-service:getDomainName() or "@"%>/gitea.ini"
start_stop_daemon_args="--user ${GITEA_USER:-gitea} --env GITEA_CUSTOM=${GITEA_CUSTOM:-/var/lib/gitea} --chdir /var/lib/gitea"
command_background=yes
pidfile="/run/yunorock/<%-service:getDomainName()%>-gitea.pid"

depend() {
        use logger dns
        need net
        after firewall
}
